# docker/systemd/secondlayer-health-alert.service
[Unit]
Description=Secondlayer prod health check with Slack alerting
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
EnvironmentFile=-/opt/secondlayer/docker/.env
ExecStart=/opt/secondlayer/docker/scripts/health-alert.sh
StandardOutput=append:/var/log/secondlayer-health-alert.log
StandardError=append:/var/log/secondlayer-health-alert.log
